# lab-knowledge-index.service
[Unit]
Description=Rebuild lab knowledge index
After=ollama.service

[Service]
Type=oneshot
User=root
ExecStart=/opt/lab-server/.venv-lab-mcp/lab-mcp/bin/python -u \
    /opt/lab-server/lab-knowledge-index.py \
    --root /opt/lab-knowledge
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lab-knowledge-index
# Allow up to 4 hours for large repos (first run may embed many files)
TimeoutStartSec=14400
