# lab-websearch.service
[Unit]
Description=Lab Web Search MCP Server (wraps SearXNG, port 3003)
Documentation=https://github.com/<lab-org>/lab-llm-server
After=network.target docker.service

[Service]
Type=simple
User=root
Group=root

ExecStart=/opt/lab-server/.venv-lab-mcp/bin/python -u \
    /opt/lab-server/lab-websearch-server.py

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lab-websearch

[Install]
WantedBy=multi-user.target
